# font-autoscale/font-autoscale.service
[Service]
Type=simple
ExecStart=/usr/local/bin/font-autoscale
PIDFile=/tmp/font-autoscale.pid

[Install]
WantedBy=multi-user.target